# etc/supervisor.cfg
[supervisord]
logfile = ${BASE_DIR}/var/log/supervisord.log
logfile_maxbytes = 10MB
logfile_backups=10
loglevel = info
pidfile = ${BASE_DIR}/var/run/supervisord.pid
childlogdir = ${BASE_DIR}/var/log
nocleanup = true
nodaemon = true

[program:zk1]
command=bash zkServer.sh start-foreground zk.cfg
environment=ZOOCFGDIR="${BASE_DIR}/var/cfg/zk1"
directory=${BASE_DIR}/comp/zookeeper/bin
process_name=zk1
redirect_stderr=true
stdout_logfile=${BASE_DIR}/var/log/supervisor/zk1.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=NONE

[program:zk2]
command=bash zkServer.sh start-foreground zk.cfg
environment=ZOOCFGDIR="${BASE_DIR}/var/cfg/zk2"
directory=${BASE_DIR}/comp/zookeeper/bin
process_name=zk2
redirect_stderr=true
stdout_logfile=${BASE_DIR}/var/log/supervisor/zk2.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=NONE

[program:zk3]
command=bash zkServer.sh start-foreground zk.cfg
environment=ZOOCFGDIR="${BASE_DIR}/var/cfg/zk3"
directory=${BASE_DIR}/comp/zookeeper/bin
process_name=zk3
redirect_stderr=true
stdout_logfile=${BASE_DIR}/var/log/supervisor/zk3.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=NONE

[program:mesos-master1]
command=mesos master --zk=zk://172.17.100.1:2181,172.17.100.2:2181,172.17.100.3:2181/mesos --quorum=2 --ip=172.17.100.4 --work_dir=${BASE_DIR}/var/data/mesos-master1
process_name=mesos-master1
redirect_stderr=true
stdout_logfile=${BASE_DIR}/var/log/supervisor/mesos-master1.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=NONE

[program:mesos-master2]
command=mesos master --zk=zk://172.17.100.1:2181,172.17.100.2:2181,172.17.100.3:2181/mesos --quorum=2 --ip=172.17.100.5 --work_dir=${BASE_DIR}/var/data/mesos-master2
process_name=mesos-master2
redirect_stderr=true
stdout_logfile=${BASE_DIR}/var/log/supervisor/mesos-master2.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=NONE

[program:mesos-master3]
command=mesos master --zk=zk://172.17.100.1:2181,172.17.100.2:2181,172.17.100.3:2181/mesos --quorum=2 --ip=172.17.100.6 --work_dir=${BASE_DIR}/var/data/mesos-master3
process_name=mesos-master3
redirect_stderr=true
stdout_logfile=${BASE_DIR}/var/log/supervisor/mesos-master3.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=NONE

[supervisorctl]
serverurl = unix:///${BASE_DIR}/var/run/supervisor.sock

[unix_http_server]
file = ${BASE_DIR}/var/run/supervisor.sock
chmod = 0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
